# bitrix/crontab/crontab.cfg

33 23 * * * /usr/local/bin/php -f /var/www/poligon/data/www/poligon.info/bitrix/php_interface/include/catalog_import/cron_frame.php 20 >/var/www/poligon/data/www/poligon.info/bitrix/php_interface/include/catalog_import/logs/20.txt
